# deploy/polybot.service
[Unit]
Description=Polymarket 5m Bot
After=network.target

[Service]
Type=simple
User=root
# CUSTOMIZE: Replace /path/to/oracel with the actual repository directory
WorkingDirectory=/path/to/oracel
# CUSTOMIZE: Replace /path/to/oracel with the actual repository directory
ExecStart=/path/to/oracel/target/release/polybot
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
