# deploy/polybot.service
[Unit]
Description=Polymarket 5m Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/polymarket-5m-bot
ExecStart=/root/polymarket-5m-bot/target/release/polybot
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
